# root/etc/s6-overlay/s6-rc.d/init-sickgear-config/run
#!/usr/bin/with-contenv bash

#make our folders and links
mkdir -p \
    /config

# If needed add dafault config
if [[ ! -f /config/config.ini ]]; then
    cp /defaults/config.ini /config/config.ini
fi

# permissions
chown -R abc:abc \
    /app \
    /config
